# python_project/videoproject/gunicorn.service
[Unit]
Description=gunicorn daemon
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/home/fantom/videoproject
ExecStart=/path/to/your/virtualenv/bin/gunicorn --workers 3 --bind unix:/path/to/your/project/yourprojectname.sock yourprojectname.wsgi:application

[Install]
WantedBy=multi-user.target